# data/rygel.service.in
[Unit]
Description=Rygel DLNA/UPnP server

[Service]
Type=dbus
BusName=org.gnome.Rygel1
ExecStart=@bindir@/rygel
Restart=on-failure

[Install]
WantedBy=default.target
